# config/systemd/user/settheme-on-wallpaper.path
[Unit]
Description=Watch swww cache for wallpaper change

[Path]
PathChanged=%h/.cache/swww
Unit=settheme-on-wallpaper.service
