# config/systemd/user/theme.target
[Unit]
Description=User Theme Stack Ready
Wants=graphical-session.target
After=graphical-session.target
